# deploy/systemd/click-send-restore-drill.timer.example
[Unit]
Description=Run Click Send Shop restore drill daily

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
