# deploy/systemd/click-send-binlog-sync.timer.example
[Unit]
Description=Run Click Send Shop binlog sync every minute

[Timer]
OnCalendar=*:*:00
AccuracySec=10s
Persistent=true

[Install]
WantedBy=timers.target
